# src/core/etc/udev/rules.d/70-net.rules
SUBSYSTEM=="net", ACTION=="add", RUN+="/sbin/udhcpc -b -i %k -x hostname:$(/bin/hostname) -p /var/run/udhcpc.%k.pid >/dev/null 2>&1 &"
